# python_apps/airtime_analyzer/install/systemd/libretime-analyzer.service
[Unit]
Description=LibreTime Media Analyzer Service

[Service]
ExecStart=/usr/bin/airtime_analyzer
User=libretime-analyzer
Group=libretime-analyzer
Restart=always

[Install]
WantedBy=multi-user.target
