# installer/systemd/libretime.target
[Unit]
Description=LibreTime Services
Wants=libretime-analyzer.service
Wants=libretime-api.service
Wants=libretime-liquidsoap.service
Wants=libretime-playout.service
Wants=libretime-worker.service

[Install]
WantedBy=multi-user.target
