# python_apps/libretime_watch/install/systemd/libretime-watch.service
[Unit]
Description=LibreTime Watch Folder Service

[Service]
ExecStart=/usr/bin/libretime-watch
Restart=always

[Install]
WantedBy=multi-user.target
