# python_apps/airtime_analyzer/install/upstart/airtime_analyzer.conf
description "Airtime Analyzer"
author "help@sourcefabric.org"

start on runlevel [2345]
stop on runlevel [!2345]

expect daemon
respawn

exec sudo -u www-data airtime_analyzer --daemon

