# analyzer/install/systemd/libretime-analyzer.service
[Unit]
Description=LibreTime Media Analyzer Service
PartOf=libretime.target

[Service]
Environment=LIBRETIME_LOG_FILEPATH=@@LOG_DIR@@/analyzer.log
Environment=LIBRETIME_CONFIG_FILEPATH=@@CONFIG_FILEPATH@@
WorkingDirectory=@@WORKING_DIR@@/analyzer

ExecStart=/usr/local/bin/libretime-analyzer
User=libretime-analyzer
Group=libretime-analyzer
Restart=always

[Install]
WantedBy=multi-user.target
