# playout/install/systemd/libretime-liquidsoap.service
[Unit]
Description=LibreTime Liquidsoap Service
PartOf=libretime.target

[Service]
Environment=LIBRETIME_LOG_FILEPATH=@@LOG_DIR@@/liquidsoap.log
Environment=LIBRETIME_CONFIG_FILEPATH=@@CONFIG_FILEPATH@@
WorkingDirectory=@@WORKING_DIR@@/playout

ExecStart=/usr/local/bin/libretime-liquidsoap
User=libretime-playout
Group=libretime-playout
Restart=always

[Install]
WantedBy=multi-user.target
