# playout/install/systemd/libretime-playout.service
[Unit]
Description=Libretime Playout Service
After=network-online.target

[Service]
ExecStart=/usr/local/bin/libretime-playout
User=libretime-pypo
Group=libretime-pypo
Restart=always

[Install]
WantedBy=multi-user.target
